# xrp_matcher/xrp_matcher.service
[Unit]
Description=XRP Matcher Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/full/path/to/xrp_matcher
ExecStart=/usr/bin/python3 /full/path/to/xrp_matcher/src/main.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=xrp_matcher
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
